# Assignement_4_MP3_Player_v1/mp3_overlay/etc/profile.d/Profile_script.sh
#!/bin/sh


PS1="\d \t Mp3_shell:\w> "

echo "Welcome to BuildRoot MP3 Player!"


